# voyant/supervisord.conf
[supervisord]
nodaemon = true

[program:tomcat]
directory=/usr/local/tomcat
startsecs=5
stopsignal=INT
command = /usr/local/tomcat/bin/catalina.sh run
redirect_stderr = true

